# statusboard-gosquared.conf
# Upstart script
# Mostly taken from http://kvz.io/blog/2009/12/15/run-nodejs-as-a-service-on-ubuntu-karmic/
description "node.js server for StatusBoard GoSquared endpoint"
author      "Dan Palmer - http://danpalmer.me/"

start on startup
stop on shutdown

respawn
respawn limit 20 5

limit nofile 32768 32768

script
  export HOME="/root"
  set -e
  exec sudo -u www-data /usr/bin/node /home/dan/repositories/statusboard-gosquared/api.js >> /var/log/node-statusboard-gosquared.log 2>&1
end script

post-start script
  # Optional - script for notifications of restart
end script